# rc.local
#!/bin/sh
ln -s /data/etc/dbus-serialbattery /opt/victronenergy/dbus-serialbattery
ln -s /data/etc/dbus-serialbattery /opt/victronenergy/service/dbus-serialbattery
