# rc.local
#!/bin/bash
sh /data/etc/dbus-serialbattery/reinstalllocal.sh
