# rc.local
#!/bin/sh
ln -s /data/etc/dbus-serialbattery/service /service/dbus-serialbattery
