# rc.local
#!/bin/sh
ln -s /data/etc/dbus-serialbattery /opt/victronenergy/dbus-serialbattery
