# rc.local
#!/bin/bash
exec /data/etc/dbus-serialbattery/reinstalllocal.sh
