# etc/dbus-serialbattery/service/log/run
#!/bin/sh
exec 2>&1
exec multilog t s2500000 n4 /var/log/dbus-serialbattery.TTY
